# supervisor-daemon.conf
[program:robot]
command=/home/pi/robot/bluebot.sh
user=root
stderr_logfile=/var/log/robot.err.log
stdout_logfile=/var/log/robot.out.log
